# files/systemd/system/deploy@.path
[Unit]
Description=Deploy %i when deployment files change!

[Path]
PathChanged=/etc/deploy/%i/
TriggerLimitIntervalSec=60
TriggerLimitBurst=60
Unit=deploy@%i.service
