# files/systemd/user/deploy@.path
[Unit]
Description=Deploy %i when deployment files change!

[Path]
PathChanged=%h/.config/deploy/%i/
TriggerLimitIntervalSec=60
TriggerLimitBurst=60
Unit=deploy@%i.service
